# systemd/greenwave.socket
[Unit]
Description=Greenwave Socket

[Socket]
ListenStream=/run/greenwave/socket

[Install]
WantedBy=sockets.target
